# roles/mail_server/templates/postfix/postfix.cron
18 21 * * fri root mysqldump service_postfix --skip-comments --skip-extended-insert --result-file=/etc/postfix/service_postfix.mysql.dump

