# roles/common/templates/systemd/system/sshd.service.d/ansible.conf
[Service]
Nice=-10
IOSchedulingClass=1
IOSchedulingPriority=0
