# roles/mail_server/templates/cron.d/dovecot
18 23 * * * root doveadm  expunge -A  mailbox Trash savedbefore 30d
