# roles/mail_server/templates/cron.d/postfix
18 21 * * fri root umask 077; mysqldump {{ postfix_db_name }} --skip-comments --skip-extended-insert --result-file=/etc/postfix/{{ postfix_db_name }}_dump.sql
40 21 * * fri root umask 077; mysqldump {{ roundcube_db_name }} --skip-comments --skip-extended-insert --no-create-info --ignore-table={{ roundcube_db_name }}.session --result-file=/etc/postfix/roundcube_data.sql
